# battery.service
[Unit]
Description=Battery-check

[Service]
Type=simple
ExecStart=/bin/sh -c '~/dotfiles/battery.sh'
