# etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --login-options=scekron 38400 tty1 linux
